# www/tlstunnel/files/tlstunnel.in
#!/bin/sh
#

# PROVIDE: tlstunnel
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add these lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# tlstunnel_enable (bool):	Set to NO by default.
#				Set it to YES to enable tlstunnel.
# tlstunnel_config (path):	Set to %%ETCDIR%%/config
#				by default.
# tlstunnel_user (user):	Set to 'tlstunnel' by default
# tlstunnel_group (group):	Set to 'tlstunnel' by default

. /etc/rc.subr

name=tlstunnel
rcvar="${name}_enable"
desc='TLS reverse proxy'

load_rc_config "${name}"
: ${tlstunnel_enable:=NO}
: ${tlstunnel_config:=%%ETCDIR%%/config}
: ${tlstunnel_user:=%%DAEMON_USER%%}
: ${tlstunnel_group:=%%DAEMON_GROUP%%}

command=/usr/sbin/daemon
procname="%%PREFIX%%/bin/${name}"
rundir="/var/run/${name}"
pidfile="${rundir}/${name}.pid"
command_args="-cfHS -p ${pidfile} -P ${rundir}/daemon_${name}.pid ${procname} -config ${tlstunnel_config}"
start_precmd="install -d -o ${tlstunnel_user} -g ${tlstunnel_group} -m 0700 ${rundir}"

run_rc_command "$1"
