# irc/soju/files/soju.in
#!/bin/sh
#

# PROVIDE: soju
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add these lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# soju_enable (bool):	Set to NO by default.
#			Set it to YES to enable soju.
# soju_config (path):	Set to %%PREFIX%%/etc/soju/config
#			by default.
# soju_user (user):	Set to 'soju' by default.
# soju_group (group):	Set to 'soju' by default.

. /etc/rc.subr

name=soju
rcvar="${name}_enable"
desc="Soju IRC bouncer"

load_rc_config "${name}"
: ${soju_enable:=NO}
: ${soju_config:=%%ETCDIR%%/config}
: ${soju_user:=%%USER%%}
: ${soju_group:=%%GROUP%%}

command=/usr/sbin/daemon
procname="%%PREFIX%%/bin/${name}"
rundir="/var/run/${name}"
pidfile="${rundir}/${name}.pid"
command_args="-cfHS -P ${rundir}/daemon_${name}.pid -p ${pidfile} ${procname} -config ${soju_config}"
start_precmd="install -d -o ${soju_user} -g ${soju_group} -m 0700 ${rundir}"

run_rc_command "$1"
