# irc/soju/files/soju.in
#!/bin/sh
#

# PROVIDE: soju
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add these lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# soju_enable (bool):	Set to NO by default.
#			Set it to YES to enable soju.
# soju_config (path):	Set to %%PREFIX%%/etc/soju/config
#			by default.

. /etc/rc.subr

name="soju"
rcvar=soju_enable
desc="Soju IRC bouncer"

: ${soju_enable:=NO}
: ${soju_config:=%%PREFIX%%/etc/soju/config}
: ${soju_user:=%%USER%%}
: ${soju_group:=%%GROUP%%}

command=/usr/sbin/daemon
procname="%%PREFIX%%/bin/soju"
soju_rundir="/var/run/${name}"
pidfile="${soju_rundir}/${name}.pid"
command_args="-P ${soju_rundir}/daemon_${name}.pid -p $pidfile -H -c -f -S $procname -config ${soju_config}"
start_precmd="install -d -o %%USER%% -g %%GROUP%% ${soju_rundir}"

load_rc_config $name
run_rc_command "$1"
